# deploy/systemd/threejmain-production-deploy-control.service
[Unit]
Description=threejmain manual production deploy control worker
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/home/threejmain
Environment=THREEJMAIN_SOURCE_REPO=/home/threejmain
Environment=THREEJMAIN_PROD_DIR=/home/threejmain-production
Environment=THREEJMAIN_PROD_BRANCH=master
Environment=THREEJMAIN_PROD_COMPOSE_PROJECT=threejmain-production
Environment=THREEJMAIN_PROD_STOP_PROJECT=threejmain
Environment=THREEJMAIN_PROD_DEPLOY_CONTROL_POLL_SECONDS=5
Environment=THREEJMAIN_PROD_COMMIT_REFRESH_SECONDS=60
Environment=GIT_CONFIG_COUNT=2
Environment=GIT_CONFIG_KEY_0=safe.directory
Environment=GIT_CONFIG_VALUE_0=/home/threejmain
Environment=GIT_CONFIG_KEY_1=safe.directory
Environment=GIT_CONFIG_VALUE_1=/home/threejmain-production
ExecStart=/home/threejmain/scripts/production_deploy_control_worker.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
